# rpmbuild/BUILD/check_device.service
[Unit]
Description=Device monitoring daemon that logs date, time, CPU and memory usage every INTERVAL_MINUTES minutes
After=network.target

[Service]
Type=forking
ExecStart=/usr/local/bin/check_device
Restart=always
User=nobody
Group=nobody

[Install]
WantedBy=multi-user.target

